# udev-files/99-usblamp.rules
SUBSYSTEM=="usb", DRIVER=="usbled", RUN+="/etc/udev/symlink-usblamp.sh"
